# build/mariadb/supervisord.conf
[supervisord]
user=root
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

[program:mariadbd]
user=root
command=/usr/local/bin/docker-entrypoint.sh mariadbd
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:cron]
user=root
command=/usr/sbin/cron -f
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
